# ops/systemd/liuxinyu-ai-loop-daily.timer
[Unit]
Description=Run Liuxinyu AI loop daily at 17:55 every day

[Timer]
OnCalendar=*-*-* 17:55:00
Persistent=true
Unit=liuxinyu-ai-loop-daily.service

[Install]
WantedBy=timers.target
